# contrib/systemd/cognito-backup.timer
[Unit]
Description=Run cognito-stack backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
